# curator4-799659099.conf
[program:curator4-799659099]
command=python3 /home/curator4/curator4-799659099/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-799659099/test.err.log
stdout_logfile=/home/curator4/curator4-799659099/test.out.log